# Projects/django/deploy_eb_docker/.conf/supervisor-app.conf
[program:uwsgi]
command = uwsgi --ini /etc/uwsgi/sites/app.ini

[program:nginx]
command = nginx